# thinkpad/etc/rc.local
#!/bin/sh -e

/usr/local/bin/thinkpad-pm-event.sh
exit 0
