# config/systemd/user/eyes-theme-boot.service
[Unit]
Description=Apply eyes theme matching current local time at session start
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=oneshot
ExecStart=%h/.scripts/eyes-theme auto

[Install]
WantedBy=default.target
