# config/systemd/user/eyes-theme-dark.service
[Unit]
Description=Switch eyes theme to dark variant

[Service]
Type=oneshot
ExecStart=%h/.scripts/eyes-theme dark
