# config/systemd/user/eyes-theme-dark.timer
[Unit]
Description=Daily 19:00 — switch eyes theme to dark

[Timer]
OnCalendar=*-*-* 19:00:00
AccuracySec=1min

[Install]
WantedBy=timers.target
