# config/systemd/user/eyes-theme-light.timer
[Unit]
Description=Daily 06:00 — switch eyes theme to light

[Timer]
OnCalendar=*-*-* 06:00:00
AccuracySec=1min

[Install]
WantedBy=timers.target
