# config/systemd/user/sway-session.target
[Unit]
Description=sway compositor session
BindsTo=graphical-session.target
Wants=graphical-session.target
After=graphical-session.target
