# etc/udev/rules.d/99-power-profile.rules
# Auto-switch power-profiles-daemon by AC state: battery -> power-saver, plugged -> balanced.
# ponytail: 2 reglas en vez de script/systemd. Requiere power-profiles-daemon activo.
SUBSYSTEM=="power_supply", ATTR{online}=="0", RUN+="/usr/bin/powerprofilesctl set power-saver"
SUBSYSTEM=="power_supply", ATTR{online}=="1", RUN+="/usr/bin/powerprofilesctl set balanced"
